# run.service
[Unit]
Description=Run script on startup

[Service]
Type=simple
ExecStart=/usr/bin/run

[Install]
WantedBy=multi-user.target
